# puppypi/start_node.service
[Unit]
Description=My Docker App
Requires=docker.service git_sync.service
After=docker.service git_sync.service

[Service]
Type=simple
Restart=no
User=pi
Environment="DISPLAY=:0"
ExecStart=/bin/bash /home/pi/puppypi/start_node.sh 

[Install]
WantedBy=multi-user.target
